# rustplayer.service
[Unit]
Description=RustPlayer Daemon
Documentation=https://github.com/yourusername/rustplayer
After=network.target sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rustplayer daemon
Restart=on-failure
RestartSec=5
Nice=10
IOSchedulingClass=idle

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=%h/.local/share/rustplayer

[Install]
WantedBy=default.target
